# .zsh/profile.d/30-rancher-desktop
#####################
# Settings for rancher desktop
#####################

rancher_desktop_bin_dir="${HOME}/.rd/bin"
if [[ -d "${rancher_desktop_bin_dir}" ]]; then
    export PATH=${rancher_desktop_bin_dir}:$PATH
fi
